# cookbooks/basics/templates/default/rc.local.erb
#!/bin/sh
#
# This script is managed by Chef. 
# Local modifications will be overwritten.
# Nuked from space to make sure its safe
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

BLUEPILLPATH=`/usr/bin/which bluepill`

touch /var/lock/subsys/local

<% if node['platform_version'].to_i <= 5 %>
/etc/init.d/rabbitmq-server start > /dev/null && $BLUEPILLPATH load /etc/bluepill/bigcouch.pill && sleep 10 && $BLUEPILLPATH load /etc/bluepill/whapps.pill > /dev/null && sleep 10 && $BLUEPILLPATH load /etc/bluepill/ecallmgr.pill > /dev/null && $BLUEPILLPATH load /etc/bluepill/freeswitch.pill > /dev/null
<% else %>
/sbin/start rabbitmq-server > /dev/null && $BLUEPILLPATH load /etc/bluepill/bigcouch.pill && sleep 10 && $BLUEPILLPATH load /etc/bluepill/whapps.pill > /dev/null && sleep 10 && $BLUEPILLPATH load /etc/bluepill/ecallmgr.pill > /dev/null && $BLUEPILLPATH load /etc/bluepill/freeswitch.pill > /dev/null
<% end %>